# debian/postrm
#!/bin/sh
#
# This script typically modifies links or other files associated with foo, 
# and/or removes files created by the package.
#

# Delete config directory if it exists
if [ -d /etc/hethio ]; then
	#don't remove directory if agent is installed
	if [ ! -f /etc/hethio/agent.conf ]; then
		rm -rf /etc/hethio
	else
		rm /etc/hethio/server.conf
	fi
fi

# Remove hethio at startup
update-rc.d -f hethio-server remove

# Remove Service Script
if [ -f /etc/init.d/hethio-server ]; then
	rm -rf /etc/init.d/hethio-server
fi

if [ -d /var/log/hethio ]; then
	if [ ! -f /var/log/hethio/agent.log ]; then
		rm -rf /var/log/hethio
	else
		rm /var/log/hethio/server*
	fi
fi